# deploy/vm/tequity.service
[Unit]
Description=Tequity Application Stack (Main + Admin)
Documentation=https://github.com/jitgadkari/tequity_monorepo
After=docker.service network-online.target
Requires=docker.service
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/tequity

# Load environment variables
EnvironmentFile=/opt/tequity/.env

# Pre-start: Pull latest images
ExecStartPre=/usr/bin/docker compose pull

# Start containers
ExecStart=/usr/bin/docker compose up -d --remove-orphans

# Stop containers gracefully
ExecStop=/usr/bin/docker compose down

# Reload: Pull new images and restart
ExecReload=/bin/bash -c '/usr/bin/docker compose pull && /usr/bin/docker compose up -d --remove-orphans'

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tequity

# Restart policy
Restart=on-failure
RestartSec=30

# Timeouts
TimeoutStartSec=300
TimeoutStopSec=60

[Install]
WantedBy=multi-user.target
